# scripts/kza-llm.service
[Unit]
Description=KZA LLM (Qwen3-30B-A3B Instruct, MoE Q5_K_M + flash_attn + KV Q8) llama.cpp HTTP server
Documentation=file:///home/kza/kza/CLAUDE.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/kza/kza
# OpenAI-compatible API en 127.0.0.1:8200 (mantiene puerto histórico del 72B
# para no romper el LLMRouter ni el `reasoner.http_base_url` del settings).
# chat_format=chatml requerido para Qwen3 family.
# n_gpu_layers=0 → CPU only (GPU 0/1 reservadas a voz/vLLM).
# Modelo: Qwen3-30B-A3B-Instruct-2507 — MoE 30B totales, 3B activos por token.
# El bandwidth efectivo de 270 GB/s × 3B ≈ 25-40 tok/s realistas (vs ~2 tok/s
# del 72B Q8_0 anterior). Calidad ArenaHard 91.0, supera DeepSeek-V3.
#
# Optimizaciones 2026-04-28 (Plan A bandwidth-limited tuning):
#   --flash_attn true        → attention con menos memory access (no acelera
#                              en CPU pero no daña, ayuda en context >100K)
#   --type_k 8 --type_v 8    → KV cache Q8_0, 50% menos memoria de cache
#                              (816MB vs 1632MB con F16). Útil para context largo
#   --use_mlock true         → garantiza no-swap (default ya era True, explícito)
#   --n_threads_batch 24     → match cores físicos (no SMT thrash en PP)
ExecStart=/home/kza/kza/.venv/bin/python -m llama_cpp.server \
  --model /home/kza/kza/models/Qwen3-30B-A3B-Instruct-2507-Q5_K_M/Qwen3-30B-A3B-Instruct-2507-Q5_K_M.gguf \
  --n_ctx 16384 --n_threads 24 --n_threads_batch 24 --n_batch 512 --n_gpu_layers 0 \
  --use_mlock true --flash_attn true --type_k 8 --type_v 8 \
  --chat_format chatml --host 127.0.0.1 --port 8200

Restart=on-failure
RestartSec=30s
StartLimitBurst=3
StartLimitIntervalSec=300

[Install]
WantedBy=default.target
